# root/etc/services.d/nginx/run
#!/usr/bin/with-contenv sh

/generate_self_signed_cert.sh

exec nginx
